# doc/90-msi-keyboard.rules
ATTR{idVendor}=="1771", ATTR{idProduct}=="ff00", GROUP="wheel", MODE="0777", SYMLINK+="msikeyboard", RUN+="/usr/local/bin/msi-udev"
